# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[program:php-fpm]
command=docker-php-entrypoint php-fpm
nodaemon=true
autostart=true
autorestart=true
